# dockerfiles/cron/conf.d/cron.conf
[program:crond]
command=/usr/sbin/cron -f -L 15
user=root
autostart=true
stdout_logfile=/var/log/supervisor/cron.log
stderr_logfile=/var/log/supervisor/cron-error.log
redirect_stderr=true