# supervisor/conf.d/php.conf
[program:php]
command=php-fpm
user=root
autostart=true
stdout_logfile=/var/log/supervisor/php.log
stderr_logfile=/var/log/supervisor/php-error.log
redirect_stderr=true